# .config/systemd/user/change-hypr-wallpaper.service
[Unit]
Description=Change Hyprpaper Wallpaper

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/scripts/change-wallpaper
